# bbs.conf
[program:bbs]
directory=/root/bbs
command=/usr/local/bin/gunicorn -b '0.0.0.0:5001' wsgi:app
autostart=true
autorestart=true
